# scripts/primordia-proxy.service
[Unit]
Description=Primordia Reverse Proxy
After=network.target

[Service]
Type=simple
User=primordia
WorkingDirectory=/home/exedev/primordia
EnvironmentFile=/home/primordia/.env.local
Environment=PRIMORDIA_WORKTREES_DIR=/home/primordia/primordia-worktrees
Environment=HOME=/home/primordia
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/bun /home/primordia/primordia-proxy.ts
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
