# scripts/primordia.service
[Unit]
Description=Primordia Next.js App
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/primordia-worktrees/current
EnvironmentFile=/home/exedev/primordia-worktrees/current/.env.local
Environment=HOSTNAME=0.0.0.0
Environment=PRIMORDIA_WORKTREES_DIR=/home/exedev/primordia-worktrees
ExecStart=/bin/bash -c 'BRANCH=$(git symbolic-ref --short HEAD 2>/dev/null || echo main) && export PORT=$(git config "branch.${BRANCH}.port" 2>/dev/null || echo 3001) && exec /home/exedev/.bun/bin/bun run start'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
