# scripts/primordia.service
[Unit]
Description=Primordia Next.js App
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/primordia-worktrees/current
EnvironmentFile=/home/exedev/primordia-worktrees/current/.env.local
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3001
Environment=PRIMORDIA_WORKTREES_DIR=/home/exedev/primordia-worktrees
ExecStartPre=/bin/bash -c 'printf "{\"port\":%%s}\\n" "${PORT:-3001}" > "${PRIMORDIA_WORKTREES_DIR:-/home/exedev/primordia-worktrees}/proxy-upstream.json"'
ExecStart=/home/exedev/.bun/bin/bun run start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
