# scripts/primordia-proxy.service
[Unit]
Description=Primordia Reverse Proxy
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/primordia-worktrees/current
EnvironmentFile=/home/exedev/primordia-worktrees/current/.env.local
Environment=PRIMORDIA_WORKTREES_DIR=/home/exedev/primordia-worktrees
ExecStart=/home/exedev/.bun/bin/bun run proxy
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
